# ambarella/prebuild/third-party/armv7-a-hf/systemd/etc/udev/rules.d/10-drivers.rules
# do not edit this file, it will be overwritten on update

ACTION=="remove", GOTO="drivers_end"

SUBSYSTEM=="platform", KERNELS=="ambarella-sd*", RUN+="/sbin/modprobe -b ambarella_sd"
SUBSYSTEM=="mmc_host", KERNELS=="mmc[0-9]", RUN+="/sbin/modprobe -b mmc_block"
SUBSYSTEM=="platform", KERNELS=="ambarella-input*", RUN+="/sbin/modprobe -b ambarella_general_input"

LABEL="drivers_end"
