# ambarella/prebuild/third-party/noarch/media-automount/etc/systemd/system/media-auto-mount@.service
[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/sbin/media-automount %I
