# ambarella/prebuild/third-party/noarch/media-automount/etc/udev/rules.d/11-usb-storage-auto-mount-systemd.rules
KERNEL=="sd[a-z][0-9]", GOTO="usb_start_mount"
KERNEL!="sd[a-z]", GOTO="usb_auto_mount_end"

LABEL="usb_start_mount"
ACTION=="add", IMPORT{program}="/sbin/blkid -o udev -p %N"
ACTION=="add", ENV{ID_FS_TYPE}=="", GOTO="usb_auto_mount_end"
ACTION=="add", RUN+="/usr/bin/systemctl restart media-auto-mount@%k.service"
ACTION=="remove", RUN+="/usr/bin/systemctl restart media-auto-mount@%k.service"
KERNEL=="sd[a-z]", ACTION=="change", RUN+="/usr/bin/systemctl restart media-auto-mount@%k.service"
LABEL="usb_auto_mount_end"
