# ambarella/prebuild/third-party/armv7-a-hf/connman/etc/udev/rules.d/70-connman-net-info.rules
# Get network device MAC address
SUBSYSTEM=="net", ACTION=="add", KERNEL=="eth[0-9]*",  RUN+="/bin/sh -c '/bin/mkdir -p /var/run/connman && /bin/echo %s{address} > /var/run/connman/%k.mac'"
SUBSYSTEM=="net", ACTION=="add", KERNEL=="wlan[0-9]*", RUN+="/bin/sh -c '/bin/mkdir -p /var/run/connman && /bin/echo %s{address} > /var/run/connman/%k.mac'"
