# ambarella/boards/s2lm_ironman/rootfs/default/etc/udev/rules.d/09-load-ucode.rules
# Load ucode when IAV driver is loaded

SUBSYSTEM=="ambpriv", KERNEL=="iav_s2l", RUN+="/usr/sbin/load_ucode /lib/firmware"
